# scripts/limofin.service
[Unit]
Description=LimoFin — personal bill & cash flow planner
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/limofin
ExecStart=/usr/bin/env node server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
